# trash-cli/.config/systemd/user/empty-trash.service
[Unit]
Description=Remove old files from trash

[Service]
Type=simple
ExecStart=trash-empty 30

[Install]
WantedBy=default.target
